# roles/ddns-update/files/ddns-update.service
[Unit]
Description=Update ddns

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ddns-update
