# roles/ddns-update/files/ddns-update.timer
[Unit]
Description=Update ddns IP-address

[Timer]
OnBootSec=0
OnUnitActiveSec=20min
AccuracySec=3min


[Install]
WantedBy=timers.target
