# roles/low-power/files/powertop.service
[Unit]
Description=Run powertop --auto-tune

[Service]
Type=oneshot
ExecStart=/usr/sbin/powertop --auto-tune
